# kmsgd.service
[Unit]
Description=Kernel Message Daemon
Documentation=man:kmsgd(8)
After=network.target local-fs.target
Wants=network-online.target

[Service]
Type=simple
ExecStart=/sbin/kmsgd start
ExecStop=/sbin/kmsgd stop
Restart=on-failure
RestartSec=5s
User=root
Group=root
StandardOutput=journal
StandardError=journal
Environment=PATH=/usr/local/sbin:/usr/local/bin:/usr/sbin:/usr/bin:/sbin:/bin

[Install]
WantedBy=multi-user.target
